# stage-halos-base/04-configure-wifi-sta-ap/files/restrict-preconfigured-wifi.service
[Unit]
Description=Restrict preconfigured WiFi to wlan0 for STA+AP mode
# Run before NetworkManager reads connection files
Before=NetworkManager.service
# Only if NetworkManager is installed
ConditionPathExists=/usr/sbin/NetworkManager

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restrict-preconfigured-wifi.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
